# misc/99.reload-monitor.rules
ACTION=="change", SUBSYSTEM=="drm", RUN+="/bin/su odysseus --command='systemctl --user start bspwm-reload.service'"
